# debian/logster-haproxy-graphite.cron.d
# crontab for logster-haproxy-graphite
# min        hour         day mon wday user     cmd
1          19            *   *   *    root  /usr/bin/graph-haproxy-metrics
