# debian/logster-haproxy-graphite.cron.weekly
#!/bin/bash

# remove pickle cache files
rm -f /var/tmp/haproxy_logster_*.p

