# debian/logster-haproxy-graphite.cron.weekly
#!/bin/bash

touch /var/tmp/haproxy_logster_cleancache

